# copy-configs/etc/systemd/system/remove-nologin.service
[Unit]
Description=Remove Nologin Blocker

[Service]
Type=oneshot
ExecStart=/usr/bin/remove-nologin

[Install]
WantedBy=multi-user.target
